# systemd/units/captdev.service
[Unit]
Description=Capt RS: Server
After=network.target

[Service]
Type=simple

User=captbaritone
Group=captbaritone

# Project directory (for relative paths, logs, etc.)
WorkingDirectory=/home/captbaritone/projects/capt-rs

# Rust server binary
ExecStart=/home/captbaritone/projects/capt-rs/deploy/server

# Optional env vars (uncomment / edit as needed)
# Environment=RUST_LOG=info
# Environment=RUST_BACKTRACE=1

Restart=on-failure
RestartSec=2

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
