# systemd/units/jordaneldredge.service
[Unit]
Description=jordaneldredge.com: Next.js Site
After=network.target

[Service]
Type=simple
User=captbaritone
Group=captbaritone
WorkingDirectory=/home/captbaritone/projects/jordaneldredge.com

Environment=HOME=/home/captbaritone
Environment=NODE_ENV=production
Environment=PORT=4000
Environment="PATH=/home/captbaritone/.nvm/versions/node/v20.19.6/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# You can keep using the absolute pnpm path if you like, but with PATH set you don't need to:
ExecStart=/usr/bin/env bash -lc 'exec pnpm run start'

Restart=on-failure
RestartSec=2

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
