# systemd/units/skin-database-blue.service
[Unit]
Description=Winamp Skin Museum: Blue Deploy
After=network.target

[Service]
Type=simple

# Run as your user
User=captbaritone
Group=captbaritone

# Project directory
WorkingDirectory=/home/captbaritone/projects/webamp/packages/skin-database/

# Environment
Environment=NODE_ENV=production
Environment=PORT=3001
Environment=HOME=/home/captbaritone
Environment="PATH=/home/captbaritone/.nvm/versions/node/v22.21.1/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Start Next.js via pnpm
ExecStart=/usr/bin/env bash -lc 'exec pnpm run start'

# Restart on failure
Restart=on-failure
RestartSec=2

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
