# units/apache-backup.timer
[Unit]
Description=Run Backup of apache config to S3 every day

[Timer]
OnCalendar=*-*-* 15:05
Persistent=true

[Install]
WantedBy=timers.target
