# units/skin-database-blue.service
[Unit]
Description=Winamp Skin Museum Blue Deploy
After=network.target

[Service]
Type=simple

# Run as your user
User=captbaritone
Group=captbaritone

# Project directory
WorkingDirectory=/home/captbaritone/projects/webamp/packages/skin-database/

# Start Next.js via pnpm
ExecStart=/home/captbaritone/.local/share/pnpm/pnpm run start

# Environment
Environment=NODE_ENV=production
Environment=PORT=3001
Environment=HOME=/home/captbaritone
Environment=PATH=/home/captbaritone/.nvm/versions/node/v22.21.1/bin:/home/captbaritone/.local/share/pnpm:/usr/bin:/bin

# Restart on failure
Restart=on-failure
RestartSec=2

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
